# etc/systemd/system/nbfc_service.service.in
[Unit]
Description=NoteBook FanControl service

[Service]
ExecStart=@BINDIR@/nbfc start
ExecStop=@BINDIR@/nbfc stop
Type=forking
PIDFile=@RUNSTATEDIR@/nbfc_service.pid
TimeoutStopSec=20
Restart=on-failure

[Install]
WantedBy=multi-user.target

