# etc/systemd/system/nbfc_service.service.in
[Unit]
Description=NoteBook FanControl service

[Service]
ExecStartPre=@BINDIR@/nbfc wait-for-hwmon
ExecStart=@BINDIR@/nbfc start
ExecStop=@BINDIR@/nbfc stop
Type=forking
PIDFile=/var/run/nbfc_service.pid
TimeoutStopSec=20
Restart=on-failure

[Install]
WantedBy=multi-user.target

